# services/claude-api-gateway.service
[Unit]
Description=Claude API Gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/claude-api-gateway start server
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
